# celery/crontab/worker_with_beat.sh
#!/bin/bash

celery -A tasks worker --loglevel=info --no-color --beat
